# applications/sensors-service/systemd-configs/sensor-collect.service
[Unit]
Description=Connect to STM32WB55 + SEN55 and store the data to DB
After=network.target

[Service]
Type=oneshot
User=
Group=
WorkingDirectory=/usr/local/bin/
EnvironmentFile=/etc/sensors-service/sensors-service.env
ExecStart=/usr/local/bin/sensors-service
TimeoutStartSec=7200
KillMode=none

StandardOutput=journal
StandardError=journal

ReadWritePaths=/var/log/sensors-service

[Install]
WantedBy=multi-user.target
